# .config/systemd/user/compton.service
[Unit]
Description=Compton X compositor

[Service]
ExecStart=/usr/bin/compton

[Install]
WantedBy=window-manager.target
