# .config/systemd/user/window-manager.target
[Unit]
Description=Graphical session with running window manager

